# dockerAssets/etc/supervisor/conf.d/nginx.conf
[program:nginx]
command = /usr/sbin/nginx -g "daemon off;"
stdout_logfile = /dev/stdout
redirect_stderr = true
autorestart = true
priority = 30
