# dockerAssets/etc/supervisor/conf.d/pm2.conf
[program:pm2]
command = pm2-runtime /var/www/%(ENV_PM2_CONFIG)s
stdout_logfile = /dev/stdout
redirect_stderr = true
autorestart = true
priority = 10
